# mumble-webui.service
[Unit]
Description=Mumble WebUI Server
After=network.target mumble-server.service
Requires=mumble-server.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/webui-mumble
ExecStart=/opt/webui-mumble/mumble-webui-server
Restart=always
RestartSec=5
Environment=PORT=8080
StandardOutput=journal
StandardError=journal

# Security settings
AmbientCapabilities=CAP_NET_BIND_SERVICE
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyDirectories=/
ReadWriteDirectories=/var/log/mumble-webui
ReadOnlyDirectories=/etc/letsencrypt/live/nimmerchat.xyz
ReadOnlyDirectories=/etc/letsencrypt/archive/nimmerchat.xyz

[Install]
WantedBy=multi-user.target