# mumble-webui.service
[Unit]
Description=Mumble WebUI Server
After=network.target mumble-server.service
Requires=mumble-server.service

[Service]
Type=simple
User=www-data
Group=mumble-web
WorkingDirectory=/opt/webui-mumble
ExecStart=/opt/webui-mumble/mumble-webui-server
Restart=always
RestartSec=5
Environment=PORT=8080
StandardOutput=journal
StandardError=journal

# Security settings
AmbientCapabilities=CAP_NET_BIND_SERVICE
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyDirectories=/
ReadWriteDirectories=/var/log/mumble-webui
ReadWriteDirectories=/opt/webui-mumble
ReadOnlyDirectories=/etc/letsencrypt

# SSL certificate access
SupplementaryGroups=ssl-cert

# Resource limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target